# deploy/vps/chatterly-api.service.example
[Unit]
Description=Chatterly API
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/srv/chatterly-v2/apps/api
Environment=NODE_ENV=production
EnvironmentFile=/srv/chatterly-v2/apps/api/.env
ExecStart=/usr/local/bin/bun run src/index.ts
Restart=always
RestartSec=3
KillSignal=SIGTERM
TimeoutStopSec=20

[Install]
WantedBy=multi-user.target
